# recipes-core/ramdisk-autostart/ramdisk-autostart/probe-autoboot.service
[Unit]
Description=Run autostart scripts

[Service]
Type=oneshot
ExecStart = /usr/bin/probe-autoboot.sh

[Install]
WantedBy=multi-user.target
