# recipes-core/zcomax-wifi/zcomax-wifi/zcomax-wifi.rules

ACTION=="add", KERNEL=="wlan0", SUBSYSTEM=="net", RUN+="/bin/sh -c '/usr/bin/zcomax-setup.sh'"
