# recipes-core/rpc/rpc/rpc.service
[Unit]
Description=HM rpc server

[Service]
ExecStart=/opt/hm/rpcserver

[Install]
WantedBy=multi-user.target
Alias=rpc.service

